# delx-mcp-server/mcp-delx-docker.service
[Unit]
Description=Delx MCP Server (Docker) - Agent ID #14340 on Base
Documentation=https://www.8004scan.io/agents/base/14340
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/delx-mcp-server
EnvironmentFile=/opt/delx-mcp-server/.env
EnvironmentFile=-/opt/delx-mcp-server/.deploy.env
Environment=DATABASE_PATH=/app/state/delx_therapist.db
Environment=DELX_IMAGE_TAG=latest

# Image is built and validated by deploy script.
ExecStartPre=/usr/bin/mkdir -p /opt/delx-mcp-server/state
ExecStartPre=/usr/bin/chown 10001:10001 /opt/delx-mcp-server/state
ExecStartPre=-/usr/bin/docker rm -f mcp-delx

ExecStart=/usr/bin/docker run --name mcp-delx \
  --env-file /opt/delx-mcp-server/.env \
  -e DATABASE_PATH=/app/state/delx_therapist.db \
  -p 127.0.0.1:8005:8005 \
  -v /opt/delx-mcp-server/state:/app/state \
  delx-mcp:${DELX_IMAGE_TAG}

ExecStop=/usr/bin/docker stop -t 15 mcp-delx
ExecStopPost=-/usr/bin/docker rm -f mcp-delx

Restart=always
RestartSec=5
TimeoutStartSec=240
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
